# DotNetCoreDebug/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:dotnetcoredebug]
command=dotnet /app/DotNetCoreDebug.dll